# net/keycloak-gatekeeper/files/keycloak-gatekeeper.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: keycloak_gatekeeper
# REQUIRE: NETWORKING SERVERS
# KEYWORD: shutdown
#
# Add these following line to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# keycloak_gatekeeper_enable (bool):       Set it to YES to enable keycloak on startup.
#                                          Default: NO
# keycloak_gatekeeper_flags (string):      Additional flags for the startup script.
#

. /etc/rc.subr

name=keycloak_gatekeeper
rcvar=keycloak_gatekeeper_enable
desc="OpenID Connect proxy server"

load_rc_config $name

: ${keycloak_gatekeeper_enable:=NO}
: ${keycloak_gatekeeper_domain:=NO}
: ${keycloak_gatekeeper_user:=%%USER%%}
: ${keycloak_gatekeeper_group:=%%GROUP%%}
: ${keycloak_gatekeeper_flags:=""}

pidfile=%%RUN_DIR%%/%%PORTNAME%%.pid
procname=%%PREFIX%%/bin/%%PORTNAME%%
command=/usr/sbin/daemon
command_args="-u ${keycloak_gatekeeper_user} -o %%LOG_DIR%%/%%PORTNAME%%.out -t %%PORTNAME%% -p %%RUN_DIR%%/%%PORTNAME%%.pid"

start_cmd="keycloak_gatekeeper_start"

keycloak_gatekeeper_start()
{
	echo "Starting %%PORTNAME%%."
	${command} ${command_args} \
		%%PREFIX%%/bin/keycloak-gatekeeper \
		--config %%PREFIX%%/etc/keycloak-gatekeeper.yaml \
		${keycloak_gatekeeper_flags}
}

run_rc_command "$1"
